# server/examples/crontab
0 * * * * /<PATH>/portal_calendar/server/cron.sh